# crontab
*/5 * * * * sudo /usr/bin/docker exec -u www-data nextcloud php -f /var/www/html/cron.php